# ops/backup/systemd/minecraft-backup-upload.service
[Unit]
Description=Cold-backup and upload Minecraft after a 30-minute quiet period
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
Environment=HOME=/root
Environment=PATH=/data/minecraft/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
LoadCredential=rcon_password:/etc/minecraft-backup/rcon-password
RuntimeDirectory=minecraft-backup
RuntimeDirectoryMode=0700
ExecStart=/usr/local/libexec/minecraft-backup/minecraft-backup-upload.sh
ExecStopPost=-/usr/local/libexec/minecraft-backup/minecraft-backup-recover.sh
TimeoutStartSec=6h
TimeoutStopSec=5min

Nice=15
IOSchedulingClass=idle
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/data/minecraft/hello-new-generation
ReadWritePaths=/data/minecraft/backups/hello-new-generation
ReadWritePaths=/data/minecraft/.config/bdpan

[Install]
WantedBy=multi-user.target
